# exam_prep/roles/bind_exporter/templates/bind_exporter.service.j2
[Unit]
Description=Bind Exporter

[Service]
User={{ bind_exp_user }}
ExecStart=/usr/local/bin/bind_exporter

[Install]
WantedBy=multi-user.target
