# Lab7/roles/apache_exporter/templates/apache_exporter.service.j2
[Unit]
Description=Apache Exporter

[Service]
User={{ apache_user }}
ExecStart=/usr/local/bin/apache_exporter $OPTIONS

[Install]
WantedBy=multi-user.target
